# core/systemd/user/no1r_reset_budget.timer
[Unit]
Description=Daily reset of no1r daily budget

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
